# vaultwarden_dev_build/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Vaultwarden
# Runs the nginx daemon
# ==============================================================================
bashio::net.wait_for 30
bashio::log.info "Starting nginx..."
exec nginx
